# sensorNode/mma8451.service
[Unit]
After=multi-user.target

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/Production/Run
StandardOutput=inherit
StandardError=inherit
Restart=no
Type=idle
ExecStart=/home/pi/Production/Run/startACQ.sh

[Install]
WantedBy=default.target
